# crontab
# Crontab for mail-mcp scheduler
#
# Copyright 2025 The Apache Software Foundation
# Licensed under the Apache License, Version 2.0
#
# Update the current month's mbox files and re-index hourly
# Runs at minute 0 of every hour
# Updates all lists configured in MAIL_MCP_MAILING_LISTS
0 * * * * /usr/local/bin/update-current-month --all
